# rc.local
#!/usr/bin/env bash
/usr/bin/findmnt /media/bd || exit 1
/usr/bin/findmnt /media/sd || exit 1
/usr/bin/findmnt /media/tv || exit 1
/usr/bin/findmnt /media/hd || exit 1
/etc/init.d/jellyfin start
